# examples/service/EchoService/echoservice.service.in
[Unit]
Description=QtService Example Echo Service
Documentation=https://github.com/Skycoder42/QtService
After=network-online.target echoservice.socket

[Service]
Type=notify
NotifyAccess=exec
ExecStart=$${target.path}/$$TARGET --backend systemd
ExecReload=$${target.path}/$$TARGET --backend systemd reload
ExecStop=$${target.path}/$$TARGET --backend systemd stop
WatchdogSec=10
Restart=on-abnormal
RuntimeDirectory=$$TARGET

[Install]
#WantedBy=multi-user.target
WantedBy=default.target
