# ProjectTemplate/service.service.in
[Unit]
Description=%{ProjectName} Service
Documentation=https://github.com/Skycoder42/QtService
@if '%{SocketPort}'
After=network-online.target %{SvcSystemdSocketName}
@endif

[Service]
Type=notify
NotifyAccess=exec
ExecStart=$${target.path}/$$TARGET --backend systemd
ExecReload=$${target.path}/$$TARGET --backend systemd reload
ExecStop=$${target.path}/$$TARGET --backend systemd stop
#WatchdogSec=10
Restart=on-abnormal
RuntimeDirectory=$$TARGET

[Install]
# Use the following for a system service
#WantedBy=multi-user.target
# Use the following for a user service
WantedBy=default.target
