# ProjectTemplate/service.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>%{SvcBundleName}</string>
	<key>ProgramArguments</key>
	<array>
		<string>$${target.path}/$$TARGET</string>
		<string>--backend</string>
		<string>launchd</string>
	</array>
	<key>RunAtLoad</key>
	<boolean>false</boolean>
@if '%{SocketPort}'
	<key>Sockets</key>
	<dict>
		<key>Listeners</key>
		<dict>
			<key>SockServiceName</key>
			<string>%{SocketPort}</string>
		</dict>
	</dict>
@endif
</dict>
</plist>
